# services/ping-monitor.service
[Unit]
Description=ISP Ping Monitor (latency + outage detection)
After=network-online.target influxdb.service
Wants=network-online.target

[Service]
Type=simple
User=__PI_USER__
WorkingDirectory=__MONITORING_DIR__
ExecStart=__MONITORING_DIR__/venv/bin/python isp/ping_monitor.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ping-monitor

[Install]
WantedBy=multi-user.target
